# images/devbox/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

[program:portpilot]
command=portpilot --listen :7681 --default 7682
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:ttyd]
command=ttyd --port 7682 --writable -t reconnect=2 tmux new-session -A -s main %(ENV_DEVBOX_APP)s
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[eventlistener:ttyd_exit]
command=/usr/local/bin/ttyd-exit-handler.sh
events=PROCESS_STATE_EXITED
