# images/devbox-agent/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

[program:portpilot]
command=portpilot --listen :7681 --default 9999
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:agent-bridge]
command=agent-bridge --listen :%(ENV_AGENT_BRIDGE_PORT)s --agent %(ENV_AGENT_BRIDGE_AGENT)s --session %(ENV_AGENT_BRIDGE_SESSION)s --db /workspace/.agent-bridge.db
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
