# plugins/agent-ollama/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ollama]
command=/bin/ollama serve
autostart=true
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=OLLAMA_HOST="0.0.0.0:11434",OLLAMA_MODELS="/data/models"

[program:sidecar]
command=/usr/local/bin/agent-ollama
autostart=true
autorestart=true
startretries=10
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
